# ntfy/rootfs/etc/services.d/example/run
#!/command/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

declare log_level
declare certfile
declare keyfile
declare base_url

bashio::config.require.ssl

# Find the matching log level
if bashio::config.has_value 'log_level'; then
    case "$(bashio::string.lower "$(bashio::config 'log_level')")" in
        all|trace)
            log_level="trace"
            ;;
        debug)
            log_level="debug"
            ;;
        info|notice)
            log_level="info"
            ;;
        warning)
            log_level="warn"
            ;;
        error|fatal)
            log_level="error"
            ;;
        off)
            log_level="off"
            ;;
    esac

    export LOG_LEVEL="${log_level}"
fi

## Print the message the user supplied, defaults to "Hello World..."
# bashio::log.info "${message:="Hello World..."}"

if bashio::config.true 'ssl'; then
  certfile=$(bashio::config 'certfile')
  keyfile=$(bashio::config 'keyfile')

  export NTFY_LISTEN_HTTPS=:7612
  export NTFY_LISTEN_HTTP=""
  export NTFY_CERT_FILE="/ssl/${certfile}"
  export NTFY_KEY_FILE="/ssl/${keyfile}"
else
  export NTFY_LISTEN_HTTP=:7612
fi

# export NTFY_BEHIND_PROXY=true
# export NTFY_WEB_ROOT=disable

export NTFY_UPSTREAM_BASE_URL=https://ntfy.sh
export NTFY_BASE_URL="$(bashio::config 'base_url')"

DATA_DIR=/data/ntfy

export NTFY_CACHE_FILE="${DATA_DIR}/cache.db"
export NTFY_AUTH_FILE="${DATA_DIR}/auth.db"
export NTFY_ATTACHMENT_CACHE_DIR="${DATA_DIR}/attachments"
export NTFY_ENABLE_LOGIN=true
export NTFY_AUTH_DEFAULT_ACCESS=deny-all

mkdir -p $NTFY_ATTACHMENT_CACHE_DIR

bashio::log.info Starting ntfy

bashio::log.info "certfile: $(bashio::config 'certfile')"
bashio::log.info "keyfile: $(bashio::config 'keyfile')"

bashio::log.info "NTFY_UPSTREAM_BASE_URL: ${NTFY_UPSTREAM_BASE_URL}"
bashio::log.info "NTFY_BASE_URL: ${NTFY_BASE_URL}"

bashio::log.info "output env: $(env)"

## Run your program
exec /usr/bin/ntfy serve
