# ntfy/rootfs/etc/services.d/example/run
#!/command/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

declare log_level
declare certfile
declare keyfile

bashio::config.require.ssl

# Find the matching log level
if bashio::config.has_value 'log_level'; then
    case "$(bashio::string.lower "$(bashio::config 'log_level')")" in
        all|trace)
            log_level="trace"
            ;;
        debug)
            log_level="debug"
            ;;
        info|notice)
            log_level="info"
            ;;
        warning)
            log_level="warn"
            ;;
        error|fatal)
            log_level="error"
            ;;
        off)
            log_level="off"
            ;;
    esac

    export LOG_LEVEL="${log_level}"
fi

## Print the message the user supplied, defaults to "Hello World..."
# bashio::log.info "${message:="Hello World..."}"

if bashio::config.true 'ssl'; then
  certfile=$(bashio::config 'certfile')
  keyfile=$(bashio::config 'keyfile')

  export NTFY_LISTEN_HTTPS=:7612
  export NTFY_LISTEN_HTTP=""
  export NTFY_CERT_FILE="${certfile}"
  export NTFY_KEY_FILE="${keyfile}"
else
  export NTFY_LISTEN_HTTP=:7612
fi

# export NTFY_BEHIND_PROXY=true
# export NTFY_WEB_ROOT=disable

bashio::log.info Starting ntfy

bashio::log.debug "certfile: $(bashio::config 'certfile')"
bashio::log.debug "keyfile: $(bashio::config 'keyfile')"

bashio::log.debug "output env: $(env)"

## Run your program
exec /usr/bin/ntfy serve
